# deployment/systemd/guard.service
[Unit]
Description=guard server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=guard
Group=guard
WorkingDirectory=/var/lib/guard
EnvironmentFile=-/etc/default/guard
# The daemon spawns ssh itself, so `guard run ssh --hostkey accept-new` needs a
# writable HOME/.ssh for ssh to record newly accepted host keys. Pin HOME to the
# state directory and have systemd pre-create .ssh (owned by the guard user).
Environment=HOME=/var/lib/guard
RuntimeDirectory=guard
RuntimeDirectoryMode=0755
StateDirectory=guard guard/.ssh
ExecStart=/usr/local/bin/guard server start --socket /run/guard/guard.sock --state-db /var/lib/guard/state.db
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/run/guard /var/lib/guard

[Install]
WantedBy=multi-user.target
