# deployment/systemd/guard-exec-as-caller.service
[Unit]
Description=guard server (execute approved commands as caller)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/var/lib/guard
Environment=HOME=/var/lib/guard
EnvironmentFile=-/etc/default/guard
RuntimeDirectory=guard
RuntimeDirectoryMode=0755
StateDirectory=guard
StateDirectoryMode=0700
ExecStartPre=/usr/bin/install -d -m 0700 /var/lib/guard/.ssh
ExecStart=/usr/local/bin/guard server start --socket /run/guard/guard.sock --users 1000 --exec-as-caller
Restart=on-failure
RestartSec=5
PrivateTmp=true
ProtectSystem=strict
ProtectHome=false
# Add any non-/home caller home roots used on this host.
ReadWritePaths=/run/guard /var/lib/guard /home

[Install]
WantedBy=multi-user.target
